# distribution/debian/romarr.service
# This file is owned by the romarr package, DO NOT MODIFY MANUALLY
# Instead use 'dpkg-reconfigure -plow romarr' to modify User/Group/UMask/-data
# Or use systemd built-in override functionality using 'systemctl edit romarr'
[Unit]
Description=Romarr Daemon
After=network.target

[Service]
User=romarr
Group=romarr
UMask=002

Type=simple
ExecStart=/opt/Romarr/Romarr -nobrowser -data=/var/lib/romarr
TimeoutStopSec=20
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
